# rpm/src/resources/systemd/hdp_04.timer
[Unit]
Description=Runs hdp_04 automatically once per hour.

[Timer]
OnCalendar=hourly
Persistent=true
